# deploy/systemd/tradingagents-engine-daily.service
# TradingAgents engine — daily scheduled run.
#
# Per spec 250 FR-033: triggered by tradingagents-engine-daily.timer at
# OnCalendar=Mon..Fri 17:00 America/New_York (after market close ET).
#
# Install path: /etc/systemd/system/tradingagents-engine-daily.service
#
# Setup:
#   sudo cp deploy/systemd/tradingagents-engine-daily.service /etc/systemd/system/
#   sudo cp deploy/systemd/tradingagents-engine-daily.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tradingagents-engine-daily.timer
#
# Manual fire (for testing):
#   sudo systemctl start tradingagents-engine-daily.service
#
# Tail logs:
#   journalctl -u tradingagents-engine-daily.service -f

[Unit]
Description=TradingAgents engine — daily watchlist run
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot

# Run as the `agent` user (uid 1000) — matches the existing agent-harness-v2
# install pattern on the VPS. SSH-verified 2026-05-10.
User=agent
Group=agent

# Working directory = the TradingAgents checkout under the agent home,
# matching /home/agent/agent-harness-v2 sibling layout (SSH-verified).
WorkingDirectory=/home/agent/tradingagents

# API keys via systemd-creds. Cred names + path match the existing
# agent-harness-v2 credstore convention (verified at /etc/credstore/):
#   anthropic-key.encrypted, exa-key.encrypted
# Provision once: sudo systemd-creds encrypt - /etc/credstore/anthropic-key.encrypted
# then paste sk-ant-... and Ctrl-D.
LoadCredentialEncrypted=anthropic-key:/etc/credstore/anthropic-key.encrypted
LoadCredentialEncrypted=exa-key:/etc/credstore/exa-key.encrypted

# Python venv path under the agent home.
Environment=PATH=/home/agent/tradingagents/.venv/bin:/usr/local/bin:/usr/bin:/bin

# Engine command. Reads ANTHROPIC_API_KEY + EXA_API_KEY from the
# CREDENTIALS_DIRECTORY exposed by systemd-creds (note: framework code
# expects ANTHROPIC_API_KEY but cred is stored under name `anthropic-key`).
ExecStart=/bin/sh -c '\
    export ANTHROPIC_API_KEY="$(cat ${CREDENTIALS_DIRECTORY}/anthropic-key)" && \
    export EXA_API_KEY="$(cat ${CREDENTIALS_DIRECTORY}/exa-key)" && \
    /home/agent/tradingagents/.venv/bin/python -m tradingagents.engine run \
        --watchlist /home/agent/tradingagents/data/watchlists/tech_weighted.txt'

# Don't restart on failure — daily timer re-fires tomorrow.
Restart=no

# 4h timeout — full 25-ticker run takes ~3-4h at 8-10 min/ticker.
TimeoutStartSec=4h

# Logs go to journald.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
